# personal_api/personal-api.service
[Unit]
Description=Personal Website API Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/personal_api
Environment="PATH=/var/www/personal_api/venv/bin"
EnvironmentFile=/var/www/personal_api/.env
ExecStart=/var/www/personal_api/venv/bin/gunicorn main:app \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 0.0.0.0:8000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --access-logfile /var/log/nginx/ps-api_access.log \
    --error-logfile /var/log/nginx/ps-api_error.log \
    --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
